# roles/dzordzu.ddf.cleaner/templates/cleaner.service.j2

[Unit]
Description=Clean home from time to time

[Service]
Type=oneshot
ExecStart=~/.ddf/bin/ddf_home_clean

[Install]
WantedBy=default.target

# vim: ft=systemd
